# ksm
#!/sbin/runscript

readonly ksm_dir=/sys/kernel/mm/ksm
readonly ksm_pages_shared=$ksm_dir/pages_shared
readonly ksm_pages_sharing=$ksm_dir/pages_sharing
readonly ksm_run=$ksm_dir/run
readonly not_supported="Kernel does not support Kernel Samepage Merging"
readonly page_size=$(getconf PAGESIZE)

btoh() {
    awk '{
        if ($1 == 1) {
            print "1 byte"
        } else {
            split("bytes KiB MiB GiB TiB PiB EiB ZiB YiB", suffix)
            for (power = 1; $1 > 1024; power++) {
                $1 /= 1024
            }
            print $1" "suffix[power]
        }
    }' <<< "$1"
}

depend() {
    before libvirtd
}

start() {
    ebegin "Starting KSM (Kernel Samepage Merging)"

    if [[ ! -f "$ksm_run" ]]; then
        eend 1 "$not_supported"
        return 1
    fi

    echo 1 > "$ksm_run"
    eend $?
}

stop() {
    ebegin "Stopping KSM (Kernel Samepage Merging)"

    if [[ ! -f "$ksm_run" ]]; then
        eend 1 "$not_supported"
        return 1
    fi

    echo 0 > "$ksm_run"
    eend $?
}

status() {
    if [[ ! -f "$ksm_run" ]]; then
        eerror "$not_supported"
        return 1
    fi

    local run=$(cat $ksm_run)
    local shared=$(cat $ksm_pages_shared)
    local sharing=$(cat $ksm_pages_sharing)
    local saved=$((sharing * page_size - shared * page_size))
    if [[ "$run" -eq 1 ]]; then
        einfo "KSM (Kernel Samepage Merging) is running; $(btoh $saved) saved"
        return 0
    else
        einfo "KSM (Kernel Samepage Merging) is not running; $(btoh $saved) \
               saved"
        return 1
    fi
}

# Local Variables:
# mode: sh
# sh-basic-offset: 4
# indent-tabs-mode: nil
# End:
