# ksm
#!/sbin/runscript

readonly KSM_DIR=/sys/kernel/mm/ksm
readonly KSM_PAGES_SHARED=$KSM_DIR/pages_shared
readonly KSM_PAGES_SHARING=$KSM_DIR/pages_sharing
readonly KSM_RUN=$KSM_DIR/run
readonly NOT_SUPPORTED="Kernel does not support Kernel Samepage Merging"
readonly PAGESIZE=$(getconf PAGESIZE)

btoh() {
    awk '{
             if ($1 == 1) {
                 print "1 byte"
             } else {
                 split("bytes KiB MiB GiB TiB PiB EiB ZiB YiB", suffix)
                 for (power = 1; $1 > 1024; power++) {
                     $1 /= 1024
                 }
                 print $1" "suffix[power]
             }
         }' <<< "$1"
}

depend() {
    before libvirtd
}

start() {
    ebegin "Starting KSM (Kernel Samepage Merging)"

    if [[ ! -f "$KSM_RUN" ]]; then
        eend 1 "$NOT_SUPPORTED"
        return 1
    fi

    echo 1 > "$KSM_RUN"
    eend $?
}

stop() {
    ebegin "Stopping KSM (Kernel Samepage Merging)"

    if [[ ! -f "$KSM_RUN" ]]; then
        eend 1 "$NOT_SUPPORTED"
        return 1
    fi

    echo 0 > "$KSM_RUN"
    eend $?
}

status() {
    if [[ ! -f "$KSM_RUN" ]]; then
        eerror "$NOT_SUPPORTED"
        return 1
    fi

    local RUN=$(cat $KSM_RUN)
    local SHARED=$(cat $KSM_PAGES_SHARED)
    local SHARING=$(cat $KSM_PAGES_SHARING)
    local SAVED=$((SHARING * PAGESIZE - SHARED * PAGESIZE))
    if [[ "$RUN" -eq 1 ]]; then
        einfo "KSM (Kernel Samepage Merging) is running; $(btoh $SAVED) saved"
        return 0
    else
        einfo "KSM (Kernel Samepage Merging) is not running; $(btoh $SAVED) \
               saved"
        return 1
    fi
}
